# ops/hetzner/dragonlens-streamlit.service
[Unit]
Description=DragonLens Streamlit UI
After=network.target dragonlens-api.service

[Service]
Type=simple
User=dragonlens
Group=dragonlens
WorkingDirectory=/opt/dragonlens
EnvironmentFile=/opt/dragonlens/.env
Environment=PYTHONPATH=/opt/dragonlens/src
ExecStart=/opt/dragonlens/.venv/bin/streamlit run src/ui/app.py --server.port 8501 --server.headless true
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dragonlens-streamlit

[Install]
WantedBy=multi-user.target
