# ops/hetzner/dragonlens-backup.service
[Unit]
Description=DragonLens PostgreSQL backup

[Service]
Type=oneshot
ExecStart=/opt/dragonlens/ops/hetzner/backup.sh
